# numerology.service
[Unit]
Description=Numerology Flask App (Gunicorn)
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/numerology
EnvironmentFile=/var/www/numerology/.env
ExecStart=/var/www/numerology/venv/bin/gunicorn -w 3 -b 127.0.0.1:8000 app:app
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
